# pw-audioscreen.service
[Service]
Type=simple
ExecStart=python '.local/share/pw-audioscreen/audio_capture.py'

[Install]
WantedBy=default.target
